# conf/supervisord-mysqld.conf
[program:mysqld]
command=/tmp/start-mysqld.sh
numprocs=1
autostart=true
autorestart=true